# scripts/com.bo.watch-self.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bo.watch-self</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/hoguej/dev/bo/scripts/watch-self-daemon.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/hoguej/dev/bo</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/hoguej/dev/bo/logs/watch-self.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/hoguej/dev/bo/logs/watch-self.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/hoguej/.bun/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
